# etc/sv/ok/run
#!/bin/bash
set -e

exec 2>&1

chpst -u nobody ok
